# dashboard/app/dashboard.service
[Unit]
Description=BigRedPi Dashboard (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/srv/homelab/dashboard/app
Environment=NODE_ENV=production
EnvironmentFile=/srv/homelab/dashboard/app/.env
ExecStart=/usr/bin/node /srv/homelab/dashboard/app/server.js
Restart=on-failure
RestartSec=3

# --- Hardening ---
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/srv/homelab/dashboard/app

[Install]
WantedBy=multi-user.target
