# dashboard/dashboard.service
[Unit]
Description=BigRedPi Dashboard
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=bigred
Group=docker
WorkingDirectory=/home/bigred/bigredpi/dashboard/app
EnvironmentFile=/home/bigred/bigredpi/dashboard/app/.env
# Optional: control refresh rate in ms (default 2000)
Environment=REFRESH_MS=2000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/bigred/bigredpi/dashboard/app
AmbientCapabilities=
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
